# third/medusa/index-configuration/medusa.cron
#!/bin/sh

killall medusa-indexd
killall medusa-searchd
/usr/bin/medusa-searchd
/usr/bin/medusa-indexd


