# debathena/debathena/clusterinfo/debian/debathena-clusterinfo.cron.hourly
#!/bin/sh

[ -x /usr/sbin/save-cluster-info ] && exec /usr/sbin/save-cluster-info
