# debathena/debathena/clusterinfo/debian/debathena-clusterinfo.cron.hourly
#!/bin/sh

exec /usr/sbin/save-cluster-info
