# debathena/debathena/command-not-found/debian/debathena-command-not-found.cron.daily
#!/bin/sh

[ -x /usr/lib/whichlocker ] && /usr/bin/curl --connect-timeout 5 https://debathena.mit.edu/auto/whichlocker-database -o /usr/share/debathena-command-not-found/program-list
