# debathena/config/cluster-login-config/debian/debathena-cluster-login-config.cron.daily
#!/bin/sh

/usr/sbin/athena-root-password
