# packs/config/os/solaris/dns-server.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	$Id: dns-server.xml,v 1.2 2005-10-18 23:48:04 rbasch Exp $
-->

<service_bundle type='manifest' name='athena:dns-server'>

<service
	name='site/athena/dns-server'
	type='service'
	version='1'>

	<dependency
	    name='filesystem_minimal'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
	    name='loopback'
	    grouping='require_any'
	    restart_on='error'
	    type='service'>
		<service_fmri value='svc:/network/loopback' />
	</dependency>

	<dependency
	    name='network'
	    grouping='optional_all'
	    restart_on='error'
	    type='service'>
		<service_fmri value='svc:/milestone/network' />
	</dependency>

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<instance name='default' enabled='true' >
		<!--
			To be able to run multiple named with their own config
			file, this dependency needs to be at the instance level.
		-->
		<dependency
		    name='config_data'
		    grouping='require_all'
		    restart_on='none'
		    type='path'>
			<service_fmri value='file://localhost/etc/named.conf' />
		</dependency>

		<!--
			Make the dns client service dependent on this, to
			ensure that syslog is started after named.
		-->
		<dependent
		    name='dns_client'
		    grouping='optional_all'
		    restart_on='none'>
		    <service_fmri value='svc:/network/dns/client' />
		</dependent>
	  
		<exec_method
		    type='method'
		    name='start'
		    exec='/etc/athena/svc-named'
		    timeout_seconds='60' >
			<method_context>
			<!--
				privileges: file_dac_read, file_dac_search
				privileges are necessary for reading the
				configuration file even it is restricted by
				the file permission. sys_resource privilege
				is for setting the resource limits (eg. stack
				size)
			-->
				<method_credential
				    user='root'
				    group='root'
				    privileges='basic,!proc_session,!proc_info,!file_link_any,net_privaddr,file_dac_read,file_dac_search,sys_resource' />
			</method_context>
		</exec_method>

		<property_group name='general' type='framework'>
			<!-- manage DNS server state -->
			<propval name='action_authorization' type='astring'
			    value='solaris.smf.manage.bind' />
			<propval name='value_authorization' type='astring'
			    value='solaris.smf.manage.bind' />
		</property_group>
	</instance>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				BIND DNS server
			</loctext>
		</common_name>
		<documentation>
			<manpage title='named' section='8'
			    manpath='/usr/athena/man' />
		</documentation>
	</template>

</service>

</service_bundle>
