# debathena/debathena/counterlog/debian/debathena-counterlog.cron.daily
#!/bin/sh

[ -x /usr/bin/athena-counterlog ] && /usr/bin/athena-counterlog
