# debathena/debathena/counterlog/debian/debathena-counterlog.cron.daily
#!/bin/sh

[ -x /usr/lib/debathena-counterlog/athena-counterlog ] && /usr/lib/debathena-counterlog/athena-counterlog
