# debathena/scripts/build-server/autolivebuilder.cron
#min	hour		mday	mon	wday	command
0	5		*	*	1	$HOME/autolivebuilder

