# debathena/scripts/build-server/autodebathenify.cron
#min	hour		mday	mon	wday	command
0	0,4,8,12,16	*	*	*	$HOME/autodebathenify

