# packs/config/os/linux/athena.cron
20 4 * * *	/etc/athena/clean_tmp_areas > /dev/null
37 * * * *	/etc/athena/desync; /usr/lib/sendmail -q
45 0,12 * * *	/etc/athena/desync 39600; . /etc/athena/rc.conf; /etc/athena/gettime -g 60 -s $TIMEHUB > /dev/null
